# 20.10.23/src/appserver.conf
[program:appserver]
user=appserver
command=python3 /app/main.py
autorestart=true
stopsignal=INT
